# infra/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/app/logs

[program:prometheus]
command=/usr/local/bin/prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/tmp/prometheus --web.listen-address=:9090 --web.enable-lifecycle
priority=10
autorestart=true
startretries=3
stdout_logfile=/app/logs/prometheus.log
redirect_stderr=true

[program:grafana]
command=/usr/sbin/grafana-server --homepath=/usr/share/grafana --config=/etc/grafana/grafana.ini
priority=20
autorestart=true
startretries=3
stdout_logfile=/app/logs/grafana.log
redirect_stderr=true
environment=GF_SECURITY_ADMIN_PASSWORD="admin",GF_AUTH_ANONYMOUS_ENABLED="true",GF_AUTH_ANONYMOUS_ORG_ROLE="Viewer",GF_PATHS_PROVISIONING="/etc/grafana/provisioning"

[program:metrics_server]
command=python -c "import framework.observability; from prometheus_client import start_http_server; start_http_server(9091, addr='0.0.0.0'); import time; time.sleep(999999)"
directory=/app
priority=30
autorestart=true
startretries=3
stdout_logfile=/app/logs/metrics_server.log
redirect_stderr=true

[program:pipeline]
command=/app/scripts/run_all_pipeline.sh
directory=/app
priority=40
autorestart=false
startsecs=0
exitcodes=0,1
stdout_logfile=/app/logs/pipeline.log
redirect_stderr=true

[program:streamlit]
command=streamlit run ui/app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true
directory=/app
priority=50
autorestart=true
startretries=3
stdout_logfile=/app/logs/streamlit.log
redirect_stderr=true
